# autostart/config-autostart/repairview.desktop
[Desktop Entry]
Type=Application
Name=RepairDashBoard
Exec=/home/cafe/repaircafe/autostart/startbrowser

